# config/backup_files.list
/etc/apt/sources.list
/etc/apt/sources.list.d/
/etc/default/locale
/etc/docker/daemon.json
/etc/group
/etc/passwd
/lib/systemd/system/docker.service

/home/pax/bin/
/home/pax/.profile
/home/pax/.bashrc
/home/pax/.xshrc
/home/pax/.ssh/
/home/pax/.editorconfig
/home/pax/.gitattributes
/home/pax/.gitconfig
/home/pax/.gitignore
/home/pax/.vimrc
/home/pax/.tmux.conf
/home/pax/.npmrc
/home/pax/.m2/settings.xml
/home/pax/.gradle/init.gradle
/home/pax/.docker/config.json
